# worker_check.service
[Unit]
Description=hiveos monitoring service

[Service]
Type=oneshot
User=root
ExecStart=/home/user/scripts/worker_check.sh
