# worker_check.timer
[Unit]
Description=hiveos monitoring timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=7min

[Install]
WantedBy=timers.target
